# ops/systemd/aps-scrape.timer
[Unit]
Description=Trigger APS scrape daily at 03:00 local
Requires=aps-scrape.service

[Timer]
OnCalendar=03:00
RandomizedDelaySec=30min
Persistent=true
Unit=aps-scrape.service

[Install]
WantedBy=timers.target
